# supervisor-sshd.conf
[program:sshd]
command=/root/scripts/sshd.sh
numprocs=1
autostart=true
autorestart=true
priority=2000