# ubuntu-setup.conf
[program:ubuntu]
command=/root/scripts/ubuntu.start.sh
autostart=false
autorestart=false

[program:cron]
command=cron -f
priority=100

[program:rsyslog]
command=rsyslogd -n
priority=10